# docker/crontab
# Crontab para Laravel Scheduler
# Formato: minuto hora dia mês dia-da-semana usuário comando
#
# IMPORTANTE: Todos os comandos devem usar o wrapper script
# que configura o ambiente PHP corretamente

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Verificar pagamentos pendentes a cada 2 horas
0 */2 * * * www-data /usr/local/bin/laravel-cron.sh pagamentos:verificar-pendentes --horas=1 >> /var/log/cron.log 2>&1

# Verificar assinaturas expiradas diariamente às 2h da manhã
0 2 * * * www-data /usr/local/bin/laravel-cron.sh assinaturas:verificar-expiradas --bloquear >> /var/log/cron.log 2>&1

# Verificar documentos de habilitação vencendo/vencidos diariamente às 6h
0 6 * * * www-data /usr/local/bin/laravel-cron.sh documentos:vencimento >> /var/log/cron.log 2>&1

# Cleanup de uploads diariamente às 3h30
30 3 * * * www-data /usr/local/bin/laravel-cron.sh documentos:cleanup-processos >> /var/log/cron.log 2>&1

# Atualizar status de processos automaticamente (a cada 5 minutos)
# Nota: Este comando também está no routes/console.php executando a cada minuto via Laravel Scheduler
# No Docker, executamos a cada 5 minutos para reduzir carga
*/5 * * * * www-data /usr/local/bin/laravel-cron.sh processos:atualizar-status-automatico >> /var/log/cron.log 2>&1
